# integrations/osquery/packs/maccrab.conf
{
  "_comment": "MacCrab companion osquery pack. 12 queries that complement MacCrab's real-time detection with point-in-time posture snapshots. Every query is macOS-only and explicitly platform='darwin'. Install by referencing this file from your osquery config's packs block.",
  "platform": "darwin",
  "version": "1.0.0",
  "queries": {
    "listening_ports_non_system": {
      "query": "SELECT pid, port, protocol, address, path FROM listening_ports lp LEFT JOIN processes p ON lp.pid = p.pid WHERE lp.address != '127.0.0.1' AND lp.address != '::1' AND (p.path NOT LIKE '/usr/libexec/%' AND p.path NOT LIKE '/System/%');",
      "interval": 900,
      "description": "Non-loopback listening ports on the host, excluding Apple's usual helpers. Pairs with MacCrab network events to spot unexpected servers.",
      "snapshot": true
    },
    "unsigned_launch_agents": {
      "query": "SELECT la.path, la.name, la.program_arguments, si.authority, si.team_identifier FROM launchd la LEFT JOIN signature si ON la.program = si.path WHERE la.run_at_load = 1 AND (si.authority IS NULL OR si.authority NOT LIKE '%Apple%');",
      "interval": 3600,
      "description": "LaunchAgents/Daemons set to auto-start whose target binary is not Apple-signed. Persistence audit — cross-reference with MacCrab's persistence_guard prevention.",
      "snapshot": true
    },
    "kernel_extensions": {
      "query": "SELECT name, version, path FROM kernel_extensions;",
      "interval": 86400,
      "description": "Full kext inventory. kexts are rare post-Catalina but still meaningful when present. Daily snapshot.",
      "snapshot": true
    },
    "startup_items_audit": {
      "query": "SELECT type, path, name, username FROM startup_items;",
      "interval": 3600,
      "description": "Every startup mechanism osquery knows about (LaunchAgents, LaunchDaemons, login hooks, login items). Complements MacCrab's real-time persistence detection.",
      "snapshot": true
    },
    "processes_with_env_dyld": {
      "query": "SELECT pid, name, path, cmdline FROM processes WHERE pid IN (SELECT pid FROM process_envs WHERE key LIKE 'DYLD_%' OR key = 'LD_PRELOAD');",
      "interval": 600,
      "description": "Processes running with DYLD_* or LD_PRELOAD injection. Paranoia snapshot — MacCrab's EnvCapture already flags these at exec time, but here's a running-snapshot view.",
      "snapshot": true
    },
    "recent_package_installs": {
      "query": "SELECT package_id, version, install_time, package_filename FROM package_receipts WHERE install_time > (strftime('%s', 'now') - (86400 * 7));",
      "interval": 3600,
      "description": "All .pkg installs in the last 7 days. Feeds supply-chain review.",
      "snapshot": true
    },
    "browser_extensions_detail": {
      "query": "SELECT uid, browser_type, name, identifier, version, description, permissions, path FROM users u CROSS JOIN browser_plugins bp WHERE bp.uid = u.uid;",
      "interval": 21600,
      "description": "Browser extension inventory across all browsers osquery recognizes. Cross-reference with MacCrab's BrowserExtensionsView.",
      "snapshot": true
    },
    "setuid_outside_system": {
      "query": "SELECT path, uid, gid, mode, size, mtime FROM file WHERE directory IN ('/usr/local/bin', '/opt/homebrew/bin', '/Users/', '/tmp', '/var/tmp') AND (mode LIKE '%4000%' OR mode LIKE '%2000%');",
      "interval": 3600,
      "description": "setuid/setgid binaries outside the standard Apple paths. Privilege-escalation posture check.",
      "snapshot": true
    },
    "login_items": {
      "query": "SELECT uid, name, path, is_hidden FROM startup_items WHERE type = 'Login Item' OR type = 'LoginItem';",
      "interval": 21600,
      "description": "Per-user login items — a common persistence location that doesn't trip LaunchAgent rules.",
      "snapshot": true
    },
    "mdns_services": {
      "query": "SELECT name, type, instance, service_port FROM mdns_services;",
      "interval": 3600,
      "description": "Services announced via Bonjour. Unexpected entries sometimes indicate C2 frameworks masquerading as AirDrop / printer helpers.",
      "snapshot": true
    },
    "non_default_routes": {
      "query": "SELECT destination, gateway, interface, netmask, flags FROM routes WHERE destination != 'default' AND destination NOT LIKE '127.%' AND destination NOT LIKE 'fe80::%';",
      "interval": 1800,
      "description": "Static routes beyond the default gateway. VPN + tunnel indicators.",
      "snapshot": true
    },
    "quarantined_executables": {
      "query": "SELECT f.path, f.size, f.mtime, xattr.value AS quarantine FROM file f CROSS JOIN extended_attributes xattr ON xattr.path = f.path WHERE f.directory LIKE '/Users/%' AND f.mode LIKE '%7%' AND xattr.key = 'com.apple.quarantine';",
      "interval": 3600,
      "description": "User-space executables with a macOS quarantine xattr. Marked files that are also executable are common post-download malware droppers.",
      "snapshot": true
    }
  }
}
