# Rules/persistence/cron_job_created.yml
title: Cron Job Created by Non-Crontab Process
id: d1a2b3c4-0012-4000-a000-000000000012
status: experimental
description: Detects modification of cron job files by processes other than crontab.
author: HawkEye Community
date: 2026/03/31
references:
    - https://attack.mitre.org/techniques/T1053/003/
tags:
    - attack.persistence
    - attack.t1053.003
logsource:
    category: file_event
    product: macos
detection:
    selection:
        TargetFilename|startswith: '/private/var/at/tabs/'
    filter_crontab:
        Image|endswith: '/crontab'
    condition: selection and not filter_crontab
falsepositives:
    - System management tools writing cron entries directly
level: high
