# Rules/persistence/launchd_override_created.yml
title: LaunchDaemon Override Database Modified
id: d1a2b3c4-0216-4000-a000-000000000216
status: stable
description: Detects modification of the launchd override database. Adversaries can modify overrides to enable disabled LaunchDaemons or alter their configuration for persistence.
author: MacCrab Community
date: 2026/03/30
references:
    - https://attack.mitre.org/techniques/T1543/004/
tags:
    - attack.persistence
    - attack.t1543.004
logsource:
    category: file_event
    product: macos
detection:
    selection:
        TargetFilename|contains:
            - '/var/db/launchd.db/'
            - '/var/db/com.apple.xpc.launchd/'
        TargetFilename|endswith:
            - 'overrides.plist'
            - '.db'
    filter_system:
        Image|startswith:
            - '/System/'
            - '/usr/libexec/'
    condition: selection and not filter_system
falsepositives:
    - System preference changes
    - Legitimate service management tools
level: high
