# Rules/persistence/launch_daemon_created.yml
title: LaunchDaemon Created
id: d1a2b3c4-0011-4000-a000-000000000011
status: stable
description: Detects creation of LaunchDaemon plist files, which run as root and persist across reboots.
author: HawkEye Community
date: 2026/03/31
references:
    - https://attack.mitre.org/techniques/T1543/004/
tags:
    - attack.persistence
    - attack.privilege_escalation
    - attack.t1543.004
logsource:
    category: file_event
    product: macos
detection:
    selection:
        TargetFilename|startswith: '/Library/LaunchDaemons/'
        TargetFilename|endswith: '.plist'
    filter_apple:
        SignerType: 'apple'
    condition: selection and not filter_apple
falsepositives:
    - Legitimate software installation requiring root persistence
level: critical
