# Rules/persistence/launch_agent_created_by_unsigned.yml
title: LaunchAgent Created by Unsigned Process
id: d1a2b3c4-0010-4000-a000-000000000010
status: stable
description: Detects creation of LaunchAgent plist files by unsigned or ad-hoc signed processes.
author: HawkEye Community
date: 2026/03/31
references:
    - https://attack.mitre.org/techniques/T1543/001/
tags:
    - attack.persistence
    - attack.t1543.001
logsource:
    category: file_event
    product: macos
detection:
    selection:
        TargetFilename|contains: '/LaunchAgents/'
        TargetFilename|endswith: '.plist'
    filter_signed:
        SignerType:
            - 'apple'
            - 'appStore'
            - 'devId'
    condition: selection and not filter_signed
falsepositives:
    - Homebrew formula installation
    - Legitimate unsigned developer tools
level: high
